# py_tls/py/tlsmonitor.service
[Unit]
Description=Truth Social Monitor Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/truth-social-monitor/py_tls/py
ExecStart=/usr/bin/python3 /root/truth-social-monitor/py_tls/py/tlsexpert.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tlsmonitor
Environment=PYTHONUNBUFFERED=1

# Give the service time to start up properly
TimeoutStartSec=60

# Restart the service on failure
StartLimitBurst=10

[Install]
WantedBy=multi-user.target 